# cron/daily-backup.cron
# Daily Backup Cron Pattern
# Runs backup script every day at 3 AM
# Documentation: https://deadmanping.com/backup-dead-man-switch

0 3 * * * /path/to/backup.sh

# With dead man switch monitoring (ping must be inside script, not here)
# The backup script should ping DeadManPing after successful backup
# If backup doesn't run or fails, ping won't arrive and DeadManPing will alert
