# cron/daily-backup.cron
# Daily Backup Cron Pattern
# Runs backup script every day at 3 AM
# Documentation: https://deadmanping.com/monitor-cron-jobs

0 3 * * * /path/to/backup.sh

# With monitoring ping (ping must be inside script, not here)
# The script should ping DeadManPing after successful backup
